# inittab
::respawn:/usr/local/bin/start_syslogd
::respawn:/usr/local/bin/start_bind
::respawn:/usr/sbin/php-fpm8 --nodaemonize --fpm-config /etc/php8/php-fpm.conf > /tmp/fpm.log 2>&1

::respawn:/usr/local/bin/start_wsgi 1
::respawn:/usr/local/bin/start_wsgi 2
::respawn:/usr/local/bin/start_wsgi 3

::respawn:/usr/sbin/nginx -g "daemon off;" >/tmp/nginx.log 2>&1

::respawn:/usr/sbin/squid --foreground
